# config/systemd/user/wob-backlight.socket
[Socket]
ListenFIFO=%t/wob-backlight.sock
SocketMode=0600

[Install]
WantedBy=sockets.target
