# etc/systemd/system/update-mirrors.service
[Unit]
Description=Update Pacman mirrors

[Service]
Type=oneshot
ExecStart=update-mirrors
